# contrib/vishnu-fms.service.in
[Unit]
Description=Vishnu File Mgmt Service
After=syslog.target

[Service]
Type=simple
ExecStart=@SBIN_DIR@/fmssed @ETC_INSTALL_DIR@/vishnu/fms.cfg
User=vishnu
Group=vishnu

[Install]
WantedBy=multi-user.target