# contrib/vishnu-ims.service.in
[Unit]
Description=Vishnu Information Mgmt Service
After=syslog.target

[Service]
Type=simple
ExecStart=@SBIN_DIR@/imssed @ETC_INSTALL_DIR@/vishnu/ims.cfg
User=vishnu
Group=vishnu

[Install]
WantedBy=multi-user.target