# contrib/vishnu-tms.service.in
[Unit]
Description=Vishnu Tasks Mgmt Service
After=syslog.target

[Service]
Type=simple
ExecStart=@SBIN_DIR@/tmssed @ETC_INSTALL_DIR@/vishnu/tms.cfg
User=vishnu
Group=vishnu

[Install]
WantedBy=multi-user.target