# contrib/packaging-obs/vishnu-X.Y.Z/debian.vishnu-server.postinst
#!/bin/bash
 
set -e

services="dispatcher umssed fmssed tmssed"
for srv in; do
 update-rc.d vishnu-$srv defaults
done


